# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/hackerthon/bin/uwsgi --ini /srv/hackerthon/.config/uwsgi/uwsgi-app.ini

