# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/hackarthon/bin/uwsgi --ini /srv/hackarthon/.config/uwsgi/uwsgi-app.ini

